# deploy/etc/services.d/backend.ini
[program:backend]
command=/home/ambience/bin/backend