# deploy/etc/services.d/nlp.ini
[program:nlp]
command=/home/ambience/bin/nlp